# templates/systemd/backup.timer.j2
[Unit]
Description=Make daily backups of GitLab

[Timer]
OnCalendar=23:23

[Install]
WantedBy=timers.target
